# src/data/supervisor/daphne.conf
[program:daphne-spms]
directory=/Users/spunkmars/SM-H/django-spms/src
user=spunkmars
command=/Users/spunkmars/SM-H/django-spms/venv/bin/daphne -b 127.0.0.1 -p 8007 --proxy-headers config.asgi:application
stdout_logfile=/Users/spunkmars/SM-H/django-spms/src/runtime/log/websocket.log
stdout_logfile_maxbytes = 20MB  ; stdout 日志文件大小，默认 50MB
stdout_logfile_backups = 20
redirect_stderr=true
autostart=true
autorestart=true
