# src/data/supervisor/gunicorn.conf
[program:gunicorn-spms]
directory=/Users/spunkmars/SM-H/django-spms/src
user=spunkmars
command=/usr/local/vpython36/bin/gunicorn -c /Users/spunkmars/SM-H/django-spms/src/data/gunicorn/gunicorn.conf -e DJANGO_ENV='prod'  config.wsgi
stdout_logfile=/Users/spunkmars/SM-H/django-spms/src/runtime/log/gunicorn-spms.log
stdout_logfile_maxbytes = 20MB  ; stdout 日志文件大小，默认 50MB
stdout_logfile_backups = 20
redirect_stderr=true
autostart=true
autorestart=true
