# src/data/supervisor/uwsgi.conf
[program:uwsgi-spms]
directory=/Users/spunkmars/SM-H/django-spms/src
user=spunkmars
environment=DJANGO_ENV='prod'
command=/usr/local/vpython36/bin/uwsgi --ini /Users/spunkmars/SM-H/django-spms/src/data/uwsgi/uwsgi_httpsocket_for_mac_9091.ini
stdout_logfile=/Users/spunkmars/SM-H/django-spms/src/runtime/log/uwsgi-spms.log
stdout_logfile_maxbytes = 20MB  ; stdout 日志文件大小，默认 50MB
stdout_logfile_backups = 20
redirect_stderr=true
autostart=true
autorestart=true
